# nps/etc/s6/nginx/run
#!/bin/sh

if [[  ! -d '/tmp/'  ]]; then
	#statements
fi

exec /usr/sbin/nginx -c /etc/nginx/nginx.conf