# root/usr/local/np-cli/etc/s6/nginx/run
#!/bin/sh

if [[  ! -d '/tmp/nginx/client-body'  ]]; then mkdir -p /tmp/nginx/client-body; fi

exec /usr/sbin/nginx -c /etc/nginx/nginx.conf