# np-cli/etc/s6/nginx/run
#!/bin/sh

if [[  ! -f "/var/log/nginx/error.log"  ]]; then 
	sudo mkdir -p /var/log/nginx
	sudo touch /var/log/nginx/error.log
	sudo chown npuser:nginx -R /var/log/nginx
	sudo chmod 775 -R /var/log/nginx
fi

if [[  ! -f "/app/log/nginx/error.log"  ]]; then 
	sudo mkdir -p /app/log/nginx
	sudo touch /app/log/nginx/error.log
	sudo chown npuser:nginx -R /app/*
	sudo chmod 775 -R /app/*
fi

if [[  ! -d "/tmp/nginx/client-body"  ]]; then
	sudo mkdir -p /tmp/nginx/client-body
	sudo chown npuser:nginx -R /tmp/nginx/client-body
	sudo chmod 775 -R /tmp/nginx/client-body
fi

exec /usr/sbin/nginx -c /etc/nginx/nginx.conf
