# nps/etc/s6/php/run
#!/bin/sh
exec /usr/bin/php-fpm -c /etc/php.ini -y /etc/php/php-fpm.conf